# .crontab
00 * * * * /home/kijima/bash/sound.sh /home/kijima/bash/win_nt.wav >> /home/kijima/Dropbox/junk/2020-01-17-024105.org

30 * * * * /home/kijima/bash/sound.sh /home/kijima/bash/st.wav >> /home/kijima/Dropbox/junk/2020-01-17-024105.org

05 * * * * /home/kijima/bash/wiki_sync.sh
