# .crontab
00 * * * * /home/kijima/bash/sound.sh /home/kijima/bash/win_nt.wav

30 * * * * /home/kijima/bash/sound.sh /home/kijima/bash/st.wav

*/5 * * * * /home/kijima/bash/wiki_sync.sh
